# recipes-ankaios/sysvinit-scripts/init-basics/init-cgroup-dhcp
#!/bin/sh

### BEGIN INIT INFO
# Provides:          network-cgroups
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start udhcpc and mount cgroups
### END INIT INFO

case "$1" in
    start)
        echo "Starting network and cgroups..."

        # Mount cgroups
        if ! mountpoint -q /sys/fs/cgroup; then
            mount -t cgroup2 none /sys/fs/cgroup
        fi

        # Start udhcpc with custom script (DNS only)
        udhcpc -i eth0 -b -s /usr/share/udhcpc/dns-only -p /var/run/udhcpc.eth0.pid
        ;;

    stop)
        echo "Stopping network and cgroups..."
        killall udhcpc
        umount /sys/fs/cgroup/cpu /sys/fs/cgroup/memory /sys/fs/cgroup 2>/dev/null
        ;;

    restart)
        $0 stop
        $0 start
        ;;

    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac

exit 0
