# recipes-ankaios/rootfs-resize/files/rootfs-resize.service
[Unit]
Description=Expand Raspberry Pi root partition and filesystem on first boot
ConditionPathExists=!/var/lib/ankaios-rpi-rootfs-resize.done
After=local-fs.target
DefaultDependencies=no
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ankaios-rpi-rootfs-resize
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
